# logstash.service.conf
[program:logstash]
command=/opt/logstash/bin/logstash -f /storage/logstash.conf
