# services/logstash.service.conf
[program:logstash]
user=docker
command=/opt/logstash/bin/logstash -f /storage/logstash.conf
stdout_events_enabled = true
stderr_events_enabled = true
